# systemd/camera-ipu4-load.service
[Unit]
Description=Load Intel IPU4 camera driver stack (Surface Pro 7)
ConditionPathExists=/home/user/camera/load-ipu4.sh

[Service]
Type=oneshot
RemainAfterExit=yes
# via /bin/bash: SELinux denies init_t direct exec of user_home_t scripts
ExecStart=/bin/bash /home/user/camera/load-ipu4.sh
